# vendor/github.com/Shopify/sarama/vagrant/kafka.conf
start on started zookeeper-ZK_PORT
stop on stopping zookeeper-ZK_PORT

# Use a script instead of exec (using env stanza leaks KAFKA_HEAP_OPTS from zookeeper)
script
  sleep 2
  export KAFKA_HEAP_OPTS="-Xmx320m"
  exec /opt/kafka-KAFKAID/bin/kafka-server-start.sh /opt/kafka-KAFKAID/config/server.properties
end script
